# deployment/etc/init/derringer.conf
pre-start script

bash << "EOF"
  mkdir -p tmp/btsync
  chown -R smashcon tmp/btsync
EOF

end script

start on runlevel [2345]

stop on runlevel [016]
