# deployment/etc/init/derringer-btsync.conf
start on starting derringer
stop on stopping derringer
